# claude-study-server.service
[Unit]
Description=Claude Study Server - AI quiz/flashcard/summary generation
After=network.target

[Service]
Type=simple
User=luis
WorkingDirectory=/home/luis/proyectos/study-framework/claude-study-server
ExecStart=/usr/bin/node /home/luis/proyectos/study-framework/claude-study-server/dist/server.js
Restart=on-failure
RestartSec=5
Environment=STUDY_SERVER_PORT=3457
Environment=CLAUDE_CODE_OAUTH_TOKEN_PERSONAL=${CLAUDE_CODE_OAUTH_TOKEN_PERSONAL}
EnvironmentFile=/home/luis/proyectos/study-framework/claude-study-server/.env

[Install]
WantedBy=multi-user.target
